# ops/systemd/glitch-ads-bot.service
[Unit]
Description=Glitch Grow Ads Agent — Telegram bot + webhook receiver
After=network.target postgresql.service meta-ads-mcp.service
Wants=postgresql.service

[Service]
Type=simple
User=support
Group=support
WorkingDirectory=/home/support/glitch-grow-ads-agent
EnvironmentFile=/home/support/glitch-grow-ads-agent/.env
ExecStart=/home/support/glitch-grow-ads-agent/.venv/bin/uvicorn ads_agent.server:app --host 127.0.0.1 --port 3110
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
