# ops/systemd/glitch-amazon-sync.service
[Unit]
Description=Glitch Grow Ads Agent — Amazon nightly cache sync
After=network.target amazon-ads-mcp.service
Wants=amazon-ads-mcp.service

[Service]
Type=oneshot
User=support
Group=support
WorkingDirectory=/home/support/glitch-grow-ads-agent
EnvironmentFile=/home/support/glitch-grow-ads-agent/.env
ExecStart=/home/support/glitch-grow-ads-agent/.venv/bin/python ops/scripts/sync_amazon.py --days 7
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glitch-amazon-sync

# 30 minute timeout — Amazon reports can stall; don't let a hung query block
# the systemd unit for hours. Next timer fire retries.
TimeoutStartSec=1800
